# config/code-server/code-server.service
# Example systemd service for code-server
# This file is created automatically by the installation script
# Location: /etc/systemd/system/code-server.service

[Unit]
Description=code-server
After=network.target

[Service]
Type=simple
User=vik9541
Environment=PASSWORD=YOUR_PASSWORD_HERE
ExecStart=/usr/bin/code-server
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
